# crontab
0 * * * * python /backupd/backup.py
@reboot python /backupd/backup.py
